# src/common/etc/systemd/system/fcoe.service.d/local.conf
[Unit]
After=xcp-networkd.service openvswitch.service xs-fcoe.service
ConditionPathExistsGlob=/etc/fcoe/cfg-eth[0-9]*
