# src/common/etc/systemd/system/vm.slice
[Unit]
Description=VM blkio slice

[Slice]
BlockIOWeight=100

[Install]
WantedBy=multi-user.target
