# src/common/etc/systemd/system/vm.slice
[Unit]
Description=VM blkio slice
Before=xapi-domains.service

[Slice]
BlockIOWeight=100

[Install]
RequiredBy=multi-user.target
