# src/common/etc/systemd/system/rsyslog.service.d/local.conf
[Unit]
DefaultDependencies=no
Requires=sysinit.target systemd-journald.service
After=sysinit.target systemd-journald.service
Before=basic.target shutdown.target
Conflicts=shutdown.target
