# src/common/etc/cron.d/logrotate.cron.rpmsave
*  * * * * root /opt/xensource/bin/logrotate-xenserver
