# src/common/etc/systemd/system/chrony-wait.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/chronyc -h 127.0.0.1,::1 waitsync 120 0.1 0.0 1
