# src/common/etc/systemd/system/xenopsd-xc.service.d/local.conf
[Unit]
Requires=xenstored.service
After=xenstored.service
