# src/common/usr/lib/systemd/system/syslog.target
[Unit]
Description=System Logging
Requires=rsyslog.service
After=rsyslog.service
