# src/common/etc/systemd/system/xapi.service.d/local.conf
[Unit]
Requires=tapback.service
After=tapback.service
