# src/common/etc/systemd/system/xapi.service.d/local.conf
[Unit]
Wants=tapback.service qemuback.service
After=tapback.service qemuback.service
