# src/common/etc/profile.d/xs.sh
#!/bin/sh

PATH=/opt/xensource/bin:$PATH
export PATH
